# pam_shield.cron
#! /bin/sh

/usr/sbin/shield-purge -c /etc/security/shield.conf

