# pam_shield.cron
#! /bin/sh

[ -x /usr/sbin/shield-purge ] && /usr/sbin/shield-purge -c /etc/security/shield.conf

