# pam_shield.cron
#! /bin/sh

if [ -x /usr/sbin/shield-purge ]; then /usr/sbin/shield-purge -c /etc/security/shield.conf
fi

